# Despliegue/mkosi.extra/etc/systemd/system/emociones.service
[Unit]
Description=Servidor de emociones

[Service]
WorkingDirectory=/Servidor
ExecStart=/Servidor/arrancar.sh
Restart=always

[Install]
WantedBy=multi-user.target
